# log-capture-webhook.service
[Unit]
Description=Xray Incident Log Capture Webhook
After=network.target uptime-kuma.service

[Service]
Type=simple
User=nobody
Group=nogroup
ExecStart=/usr/bin/python3 /opt/log-capture-webhook.py
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
