# tls-push-monitor.service
[Unit]
Description=TLS Handshake Monitor (Prometheus Push)
After=network.target

[Service]
Type=simple
User=nobody
Group=nogroup
ExecStart=/usr/local/bin/tls-push-monitor.sh
Restart=always
RestartSec=5
EnvironmentFile=/etc/xray-monitor.env

[Install]
WantedBy=multi-user.target
